# sensorcommunity/dataingester/deploy/dataingester.service
[Unit]
Description=dataingester daemon program

[Service]
Type=simple
ExecStart=/usr/local/bin/dataingester \
	-c /etc/ingester/dataingester.toml
User=root
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
